# macos/launchagents/update_kevin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>update.kevin</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/joshf/scripts/run.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
    <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>11</integer>
            <key>Minute</key>
            <integer>0</integer>
    </dict>
    <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>11</integer>
            <key>Minute</key>
            <integer>0</integer>
    </dict>
    </array>
        <key>StandardOutPath</key>
        <string>/Users/joshf/kevin.log</string>
        <key>StandardErrorPath</key>
        <string>/Users/joshf/kevin.log</string>
</dict>
</plist>
